# etc/systemd/system/nmea2file.service
[Unit]
Description=Save Airmar 200wx NMEA data to file
# https://wiki.archlinux.org/index.php/Systemd#Handling_dependencies
Requires=kplex.service
After=kplex.service

[Service]
Type=simple
RemainAfterExit=No
ExecStart=/usr/sbin/nmea2file

[Install]
WantedBy=multi-user.target

